# config/_old/rsyslog/runit/rsyslog/run
#!/bin/bash
exec 2>&1

sed -i s/\$ModLoad\ imklog/#\$ModLoad\ imklog/g /etc/rsyslog.conf

cd /
exec rsyslogd -n